# services/reboot_hung_systems/reboot_hung_systems-watcher.service
[Unit]
Description=Hung System restart watcher
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart reboot_hung_systems
